# templates/init.sh.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides: bamboo
# Required-Start: $local_fs $remote_fs $network $time
# Required-Stop: $local_fs $remote_fs $network $time
# Should-Start: $syslog
# Should-Stop: $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Atlassian Bamboo Server
### END INIT INFO
# INIT Script
######################################

# Note, this came from
#   https://confluence.atlassian.com/bamboo/running-bamboo-as-a-linux-service-416056046.html
#
# It was then hacked on a bit to actually work. :-(

# Name of app ( bamboo, Confluence, etc )
APP=<%= @username %>

# Name of the user to run as
USER=<%= @username %>

# Location of application's bin directory
BASE=<%= @home %>

export JAVA_HOME=<%= @java_home %>

case "$1" in
  start)
    chown -R <%= @username %>:<%= @username %> $BASE
    echo "Setting max locked memory"
    ulimit -l unlimited
    echo "Starting $APP"
    /bin/su - $USER -c "$BASE/bin/bamboo-agent.sh start &> /dev/null"
    ;;
  stop)
    echo "Stopping $APP"
    /bin/su - $USER -c "$BASE/bin/bamboo-agent.sh stop &> /dev/null"
    echo "$APP stopped successfully"
    ;;
  restart)
    $0 stop
    sleep 5
    $0 start
    ;;
  status)
    /bin/su - $USER -c "$BASE/bin/bamboo-agent.sh status"
    ;;
  *)
    echo "Usage: /etc/init.d/$APP {start|restart|stop}"
    exit 1
    ;;
esac

exit 0
0
